# ap-config/sbc-wlan0-up.service
[Unit]
Description=wlan0 static IP + regdom for AIC8800D40 AP
Before=hostapd.service dnsmasq.service
After=network.target sys-subsystem-net-devices-wlan0.device

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=AP_IP=192.168.43.1
ExecStart=/usr/local/sbin/sbc-ap-wifi-up.sh

[Install]
WantedBy=multi-user.target
